# systemd/rotse-weatherd.service
[Unit]
Description=Weather Station Data Collector Daemon for ROTSE-IIId Namibia Telescope and Weather Station
After=network.target
Wants=network.target

[Service]
Type=simple
User=weather
Group=weather
WorkingDirectory=/opt/rotse-weatherd
ExecStart=/usr/local/bin/rotse-weatherd --config /etc/rotse-weatherd/config.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weather-collector

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /var/lib/rotse-weatherd
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
MemoryLimit=128M

[Install]
WantedBy=multi-user.target